# system/management-api/management-api.service
[Unit]
Description=Luigi Management API Server
Documentation=https://github.com/pkathmann88/luigi
# Wait for network to be fully online before starting web server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=luigi-api
Group=luigi-api
WorkingDirectory=/var/lib/luigi-api/management-api

# Create log directory before starting
# Use + prefix to run with elevated privileges and skip WorkingDirectory restriction
ExecStartPre=+/bin/mkdir -p /var/log/luigi
# Validate configuration and certificates before starting
ExecStartPre=+/var/lib/luigi-api/management-api/scripts/pre-start-check.sh

# Environment
Environment="NODE_ENV=production"
EnvironmentFile=/etc/luigi/system/management-api/.env

# Execute
ExecStart=/usr/bin/node /var/lib/luigi-api/management-api/server.js

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=management-api

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/luigi
ReadOnlyPaths=/etc/luigi/system/management-api/certs

# Allow access to audio devices for sound playback
DeviceAllow=/dev/snd/pcmC0D0p rw
DeviceAllow=/dev/snd/controlC0 rw
DeviceAllow=/dev/snd/timer r
DevicePolicy=closed

# Resource limits (Raspberry Pi Zero W constraints)
LimitNOFILE=1024
MemoryMax=200M
CPUQuota=80%

[Install]
WantedBy=multi-user.target
