# system/system-info/system-info.service
[Unit]
Description=System Info Monitor Service
Documentation=https://github.com/pkathmann88/luigi
# Wait for network to be fully online (not just initialized)
# This ensures DNS resolution and MQTT broker connectivity work
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Service runs as dedicated user luigi-sysinfo (no special privileges needed)
User=luigi-sysinfo
Group=luigi-sysinfo
ExecStartPre=/bin/mkdir -p /var/log/luigi
ExecStart=/usr/bin/python3 /usr/local/bin/system-info.py
Restart=on-failure
RestartSec=10
# Dual logging: both journalctl (automatic) and /var/log/luigi/system-info.log (traditional)
StandardOutput=append:/var/log/luigi/system-info.log
StandardError=append:/var/log/luigi/system-info.log

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
