# system/system-info/system-info.service
[Unit]
Description=System Info Monitor Service
Documentation=https://github.com/pkathmann88/luigi
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /usr/local/bin/system-info.py
Restart=on-failure
RestartSec=10
# Dual logging: both journalctl (automatic) and /var/log/system-info.log (traditional)
StandardOutput=append:/var/log/system-info.log
StandardError=append:/var/log/system-info.log

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
