# motion-detection/mario/mario.service
[Unit]
Description=Mario Motion Detection Service
Documentation=https://github.com/pkathmann88/luigi
# Wait for network to be fully online (not just initialized)
# This ensures DNS resolution and MQTT broker connectivity work
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Service runs as dedicated user luigi-mario (member of gpio group for hardware access)
# No root privileges needed - gpio group provides necessary hardware access
User=luigi-mario
Group=luigi-mario
# Working directory must be writable by service user for lgpio library
# The lgpio library creates temporary files (.lgd-nfy*) in the working directory
WorkingDirectory=/var/lib/luigi-mario
# Ensure directories exist before starting
ExecStartPre=/bin/mkdir -p /var/log/luigi
ExecStartPre=/bin/mkdir -p /var/lib/luigi-mario
ExecStartPre=/bin/chown luigi-mario:luigi-mario /var/lib/luigi-mario
ExecStart=/usr/bin/python3 /usr/local/bin/mario.py
Restart=on-failure
RestartSec=10

# Logging: Python logging module writes to /var/log/luigi/mario.log
# Systemd journal captures stdout/stderr for 'journalctl -u mario.service'
# No StandardOutput/StandardError redirect to avoid duplicate logging

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
